# deploy/spiral-dashboard.service
[Unit]
Description=Spiral Dashboard Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/spiral
Environment="PATH=/opt/spiral/venv/bin"
EnvironmentFile=/etc/default/spiral
ExecStart=/opt/spiral/venv/bin/gunicorn -c /opt/spiral/deploy/gunicorn_conf.py 'spiral.app:create_app()'
Restart=always
RestartSec=5

# Security settings
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true
ProtectHome=true
ReadWritePaths=/var/log/spiral

[Install]
WantedBy=multi-user.target
